# config/systemd/soc-baseline.service
[Unit]
Description=LLM-based SOC - Daily Baseline Digest Capture

[Service]
Type=oneshot
User=user
WorkingDirectory=/home/user/projects/llm-soc-toolkit
ExecStart=/home/user/projects/llm-soc-toolkit/scripts/capture-baseline
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/var/lib/headless-siem
ReadWritePaths=/home/user/projects/llm-soc-toolkit/soc_context/baselines
